# ut2k4server.cron
MAILTO=mail-address@yourserv.er
MAILFROM=mail-address@yourserv.er

@reboot '/home/ut2k4server/ut2k4server-2 monitor' > /dev/null 2>&1

0 4 * * 0   ut2k4server   /usr/local/bin/ut2k4server_getplayersfromlog.sh
*/10 3-5 * * *   ut2k4server   /usr/local/bin/lgsm_restart-when-needed.sh ut2k4server doit >/dev/null 2>&1

55 3 * * *  root  systemctl restart ut2k4server >/dev/null 2>&1

